# systemd/codex-telegram-maintain.service
[Unit]
Description=Telegram Codex runner maintenance
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
# Install path. See codex-telegram-bot.service for the parameterization
# pattern and /etc/default/codex-telegram-runner for the override file.
Environment=CODEX_TELEGRAM_DIR=/opt/codex-telegram-runner
EnvironmentFile=-/etc/default/codex-telegram-runner
WorkingDirectory=${CODEX_TELEGRAM_DIR}
EnvironmentFile=${CODEX_TELEGRAM_DIR}/.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=${CODEX_TELEGRAM_DIR}/.venv/bin/python ${CODEX_TELEGRAM_DIR}/scripts/auto_maintain.py --clean-threshold 100 --keep 50
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=${CODEX_TELEGRAM_DIR} /srv /home/ubuntu
