# systemd/codex-telegram-bot.service
[Unit]
Description=Telegram-controlled Codex CLI runner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
# Install path. Defaults to /opt/codex-telegram-runner; override by
# creating /etc/default/codex-telegram-runner and setting
# CODEX_TELEGRAM_DIR there. The .venv/bin/python interpreter and the
# .env secrets file are both expected under this directory.
Environment=CODEX_TELEGRAM_DIR=/opt/codex-telegram-runner
EnvironmentFile=-/etc/default/codex-telegram-runner
WorkingDirectory=${CODEX_TELEGRAM_DIR}
EnvironmentFile=${CODEX_TELEGRAM_DIR}/.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=${CODEX_TELEGRAM_DIR}/.venv/bin/python ${CODEX_TELEGRAM_DIR}/bot.py
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=${CODEX_TELEGRAM_DIR} /srv /home/ubuntu

[Install]
WantedBy=multi-user.target
