# systemd/conveyor-feishu-bot.service
[Unit]
Description=Conveyor: Feishu-controlled Codex CLI runner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
EnvironmentFile=-/etc/default/conveyor
WorkingDirectory=/opt/conveyor
EnvironmentFile=/opt/conveyor/.env
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/conveyor/.venv/bin/python /opt/conveyor/feishu_bot.py
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
# Specific writable paths:
# - /opt/conveyor: installation directory (logs, state, etc.)
# - /srv: workspace task directories
# - /home/ubuntu/.codex: memory root storage
# - /home/ubuntu/.local/share/conveyor: local configuration and token files
ReadWritePaths=/opt/conveyor /srv /home/ubuntu/.codex /home/ubuntu/.local/share/conveyor

[Install]
WantedBy=multi-user.target
